# files_8m/etc/uci-defaults/40-network
#!/bin/sh
uci batch << EOF
delete network.wan6
set network.lan.ip6assign='64'
set xmm-modem.@xmm-modem[-1].enable='0'
commit network
EOF
