# files_16m/etc/rc.local
# Put your custom commands here that should be executed once
# the system init finished. By default this file does nothing.

# MeIG SLM828 Detect
#echo '2dee 4d22 ff 2cb7 0104' > /sys/bus/usb-serial/drivers/option1/new_id
#echo '2dee 4d22 ff 2cb7 0104' > /sys/bus/usb/drivers/qmi_wwan/new_id

# GosunCN Welink GM510 Detect
#echo '0305 1415 ff' > /sys/bus/usb-serial/drivers/option1/new_id

# ThinkWill ML7820+ Detect
#echo '19d2 0532 ff' > /sys/bus/usb-serial/drivers/option1/new_id

# Yuge CLM920 ASR Chipt Detect
#echo '1286 4e3c ff' > /sys/bus/usb-serial/drivers/option1/new_id

# Fibocom NL678 GTUSBMODE 22 (modem + AT + RMNET)
#echo '2cb7 0109 ff 12d1 1406' > /sys/bus/usb-serial/drivers/option1/new_id
#echo '2cb7 0109 ff 1bc7 1261' > /sys/bus/usb/drivers/qmi_wwan/new_id

exit 0
