# files/etc/uci-defaults/15-system
#!/bin/sh

MACHINE=$(awk '/machine/{print $4}' /proc/cpuinfo)

uci batch << EOF
set system.@system[0].zonename='Europe/Moscow'
set system.@system[0].hostname=$MACHINE
EOF
