# files/etc/uci-defaults/40-network
#!/bin/sh
uci batch << EOF
delete network.wan6
set network.azp=interface
set network.azp.device='azp0'
set network.azp.proto='none'
set network.lan=interface
set network.lan.device='br-lan'
set network.lan.proto='static'
delete network.lan.ip6assign='60'
set network.wwan=interface
set network.wwan.proto='modemmanager'
set network.wwan.device='/sys/devices/platform/101c0000.ehci/usb1/1-1/1-1.3'
set network.wwan.iptype='ipv4'
set network.wwan.auth='none'
commit network
EOF
