# files/etc/uci-defaults/80-smstools3
#!/bin/sh
EXAMLE=/usr/share/luci-app-smstools3/smstools3.user
if [ -f $EXAMPLE ]; then
	cp -rp $EXAMPLE /etc
fi
