# files/etc/uci-defaults/90-firewall
#!/bin/sh

[ -f /etc/.config.lock ] && exit 0

uci batch << EOF
del_list firewall.@zone[1].network='wan6'
commit firewall
EOF

