# files_16m/etc/uci-defaults/90-firewall
#!/bin/sh

uci batch << EOF
set firewall.@zone[1].network='wan wwan azp'
commit firewall
EOF
