# files_16m/etc/uci-defaults/90-firewall
#!/bin/sh

uci batch << EOF
add_list firewall.@zone[1].network='wan'
add_list firewall.@zone[1].network='azp'
commit firewall
EOF
