# files_16m/etc/uci-defaults/zconfig
#!/bin/sh
#
LOCK="/etc/.config.lock"

! [ -f $LOCK ] && {
	echo "$LOCK" >> /etc/sysupgrade.conf
	touch $LOCK
}

exit 0
